# templates/systemd/librebooking.service.j2
{#
SPDX-FileCopyrightText: 2026 shukon (https://github.com/shukon)

SPDX-License-Identifier: AGPL-3.0-or-later
#}

[Unit]
Description=LibreBooking ({{ librebooking_identifier }})
{% for service in librebooking_systemd_required_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}
{% for service in librebooking_systemd_wanted_services_list %}
Wants={{ service }}
{% endfor %}
DefaultDependencies=no

[Service]
Type=simple
Environment="HOME={{ devture_systemd_docker_base_systemd_unit_home_path }}"
ExecStartPre={{ devture_systemd_docker_base_host_command_sh }} -c '{{ librebooking_container_cleanup_command | trim | replace("\n", " ") }}'

{#
  Running with additional capabilities is required: the container entrypoint runs sed as www-data
  and needs CAP_DAC_OVERRIDE and CAP_FOWNER to write config files owned by root.
#}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} create \
      --rm \
      --name={{ librebooking_identifier }} \
      --log-driver=none \
      --cap-drop=ALL \
      --cap-add=CAP_DAC_OVERRIDE \
      --cap-add=CAP_FOWNER \
      --cap-add=CHOWN \
      --cap-add=SETGID \
      --cap-add=SETUID \
      --network={{ librebooking_container_network }} \
      {% if librebooking_container_http_host_bind_port %}
      -p {{ librebooking_container_http_host_bind_port }}:{{ librebooking_container_http_port }} \
      {% endif %}
      --env-file={{ librebooking_base_path }}/env \
      --label-file={{ librebooking_base_path }}/labels \
      --mount type=bind,src={{ librebooking_config_path }},dst=/config \
      --mount type=bind,src={{ librebooking_uploads_images_path }},dst=/var/www/html/Web/uploads/images \
      --mount type=bind,src={{ librebooking_uploads_reservation_path }},dst=/var/www/html/Web/uploads/reservation \
      {% for volume in librebooking_container_additional_volumes %}
      --mount type={{ volume.type | default('bind' if '/' in volume.src else 'volume') }},src={{ volume.src }},dst={{ volume.dst }}{{ (',' + volume.options) if volume.options | default('') else '' }} \
      {% endfor %}
      {% for arg in librebooking_container_extra_arguments %}
      {{ arg }} \
      {% endfor %}
      {{ librebooking_container_image }}

{% for network in librebooking_container_additional_networks %}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} network connect {{ network }} {{ librebooking_identifier }}
{% endfor %}

ExecStart={{ devture_systemd_docker_base_host_command_docker }} start --attach {{ librebooking_identifier }}

ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} stop -t {{ librebooking_container_stop_grace_time_seconds }} {{ librebooking_identifier }} 2>/dev/null || true'
ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ librebooking_identifier }} 2>/dev/null || true'

Restart=always
RestartSec=30
SyslogIdentifier={{ librebooking_identifier }}

[Install]
WantedBy=multi-user.target
